# SystemD/BLOG_server.service
[Unit]
Description=BLOG server

[Service]
User=stefus
Group=stefus
ExecStart=/usr/bin/bash /opt/BLOG/setup.sh

[Install]
WantedBy=multi-user.target
